# deploy/pmo-backup.timer
[Unit]
Description=Daily InnoGreen PMO DB backup

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
